# initra.service
[Unit]
Description=Desktop Read Ahead

[Service]
ExecStart=/usr/bin/initra-desktop

[Install]
WantedBy=multi-user.target
